# rootfs/etc/systemd/system/systemd-timedated.service.d/localtime.conf
[Service]
Environment=
Environment=SYSTEMD_ETC_LOCALTIME=/etc/localtime
ReadWritePaths=/etc
