# rootfs/etc/systemd/system/haos-wipe.service.d/override.conf
[Unit]
After=
Before=
Requires=
Wants=
ConditionPathIsMountPoint=/mnt/boot
