# rootfs/etc/systemd/system/rauc.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/rauc --mount=/run/rauc/mnt --override-boot-slot=A service
